# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ---------------------------------------------------------------------------
; PostgreSQL
; ---------------------------------------------------------------------------
[program:postgresql]
command=/usr/pgsql-%(ENV_PG_VERSION)s/bin/postgres -D %(ENV_PGSQL_DATA)s
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql-err.log

; ---------------------------------------------------------------------------
; WildFly
; ---------------------------------------------------------------------------
[program:wildfly]
command=%(ENV_WILDFLY_HOME)s/bin/standalone.sh -b 0.0.0.0 -bmanagement 0.0.0.0
autostart=true
autorestart=true
priority=20
startsecs=15
stdout_logfile=/var/log/wildfly.log
stderr_logfile=/var/log/wildfly-err.log
